# supervisor.docker.conf
[program:docker]
command=/usr/local/bin/docker-entrypoint.sh
autorestart=true
